# deploy/docker/autostart/docker-openvins-side-left.service
[Unit]
Description=ROS2 device driver container
Requires=docker.service
After=docker.service

[Service]
Restart=always
ExecStart=/usr/bin/docker run --network host --privileged openvins:side_left
ExecStop=/usr/bin/docker stop -t 2 openvins:side_left

[Install]
WantedBy=default.target
